# archlinux/systemd/.config/systemd/user/bedtime.service
[Unit]
Description=Bedtime routine script

[Service]
ExecStart=%h/.local/bin/bedtime.sh

[Install]
WantedBy=timers.target
