# archlinux/systemd/.config/systemd/user/bedtime.timer
[Unit]
Description=Timer for bedtime routine

[Timer]
OnCalendar=*-*-* 21:45:00
Persistent=false

[Install]
WantedBy=timers.target
