# archlinux/systemd/.config/systemd/user/rclone-bisync.timer
[Unit]
Description=Run Rclone Bisync periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
