# archlinux/systemd/.config/systemd/user/user-notify-fail@.service
[Unit]
Description=Notify on service failure

[Service]
Type=oneshot
ExecStart=%h/.local/bin/systemd-notify-fail %i
