# archlinux/systemd/.config/systemd/user/user-notify-success@.service
[Unit]
Description=Notify on service success

[Service]
Type=oneshot
ExecStart=%h/.local/bin/systemd-notify-success %i
